# raspberry-wake.service
[Unit]
Description=A job to test the systemd scheduler

[Service]
Type=simple
ExecStart=%h/wake.sh

[Install]
WantedBy=default.target
